# build/squeezelite.ir.conf
[supervisord]
nodaemon=true

[program:squeezelite]
command=bash -c "sleep 3; /root/squeezelite -i /root/Squeezebox3 -V Master -d all=info -d ir=debug -f /squeezebox/squeezelite.log -N /squeezebox/name"
redirect_stderr=true


[program:lirc]
command=/usr/sbin/lircd -n -L /squeezebox/lircd.log /root/Squeezebox3
redirect_stderr=true

